# systemd/containerstacks-billing-dev.service
[Unit]
Description=ContainerStacks Billing Daemon (Development)
Documentation=https://github.com/skyvps360/containerstacks
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
# Development version - uses root for simplicity
User=root
Group=root
WorkingDirectory=/root/containerstacks
ExecStart=/usr/bin/npx tsx scripts/billing-daemon/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=containerstacks-billing

# Environment configuration
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/root/.local/bin"
EnvironmentFile=/root/containerstacks/.env

# Basic security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target