# systemd/containerstacks-billing.service
[Unit]
Description=ContainerStacks Billing Daemon
Documentation=https://github.com/yourusername/containerstacks
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
# Use a dedicated user for security (create with: sudo useradd -r -s /bin/false containerstacks)
# For development/testing, you can use root, but production should use a dedicated user
User=containerstacks
Group=containerstacks
WorkingDirectory=/opt/containerstacks
ExecStart=/opt/containerstacks/node_modules/.bin/tsx scripts/billing-daemon/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=containerstacks-billing

# Environment configuration
Environment="NODE_ENV=production"
EnvironmentFile=/opt/containerstacks/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/containerstacks/logs
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
